# code/profitlab_quantum/tools/scrape_economic_calendar.timer
[Unit]
Description=Daily Economic Calendar Scraper (6:00 UTC)

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
